# remotecli/agent/systemd/sehcs-agent.timer
[Unit]
Description=Run SEHCS Agent hourly

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Unit=sehcs-agent.service

[Install]
WantedBy=timers.target